# packages/wunderland/deployment/rabbithole/rabbithole.service
[Unit]
Description=Rabbit Hole Inc - Next.js App
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/app/rabbithole
ExecStart=/usr/bin/node apps/rabbithole/.next/standalone/apps/rabbithole/start.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Environment=NODE_PATH=/app/rabbithole/node_modules
EnvironmentFile=-/etc/rabbithole/env

[Install]
WantedBy=multi-user.target
